# templates/cleanup.cron
#birdnet
*/5 * * * * pi /usr/local/bin/disk_check.sh &> /dev/null
#birdnet
*/3 * * * * pi /usr/local/bin/cleanup.sh &> /dev/null
#birdnet
@reboot pi /usr/local/bin/cleanup.sh &> /dev/null
