# templates/cleanup.cron
#birdnet
0 1 * * * /usr/local/bin/cleanup.sh &> /dev/null
#birdnet
@reboot /usr/local/bin/cleanup.sh &> /dev/null
