# templates/species_updater.cron
#birdnet
*/5 * * * * /usr/local/bin/species_notifier.sh &> /dev/null
